# config/99-xcena_set_devdax_perm.rules
ACTION=="add|change", SUBSYSTEM=="dax", RUN+="/usr/local/sbin/xcena_set_devdax_perm"
